# com.openclaw.fileorganizer.stripe-poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.fileorganizer.stripe-poller</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ripley/.openclaw/workspace/file_organizer/.venv/bin/python</string>
        <string>/Users/ripley/.openclaw/workspace/file_organizer/src/core/stripe_poller.py</string>
        <string>--hours</string>
        <string>4</string>
    </array>
    
    <!-- Run every 4 hours -->
    <key>StartInterval</key>
    <integer>14400</integer>
    
    <!-- Or use StartCalendarInterval for specific times -->
    <!--
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    -->
    
    <key>WorkingDirectory</key>
    <string>/Users/ripley/.openclaw/workspace/file_organizer</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/ripley/.openclaw/workspace/file_organizer</string>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/ripley/.openclaw/logs/stripe-poller.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/ripley/.openclaw/logs/stripe-poller-error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
